# conf.d/openbox.conf
[program:openbox]
command=openbox
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
